# target/pbx-deciso/rootfs/etc/inittab
::sysinit:/sbin/syslogd -S -C512
::sysinit:/etc/rc.bootup
ttyS0::respawn:/etc/rc.initial
