# target/pbx-cp3000/rootfs/etc/inittab
#::sysinit:/etc/rc.bootup
#tty1::respawn:/etc/rc.initial
#tty2::respawn:/sbin/getty 38400 tty2
#tty3::respawn:/sbin/getty 38400 tty3
ttyBF0::respawn:/sbin/getty 57600 ttyBF0