# target/pbx-cp3000/rootfs/etc/inittab
::sysinit:/etc/rc.bootup
ttyBF0::respawn:/etc/rc.initial
